# host/systemd/sensor-system-supervisor.service
[Unit]
Description=Sensor System RS485 multi-channel supervisor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=anone
Group=anone
WorkingDirectory=/home/anone/pico-projects
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/anone/pico-projects/host/.venv/bin/python /home/anone/pico-projects/host/host_supervisor.py --config /home/anone/pico-projects/host/system_config.json
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
